# drumroll-service_1.0.2-20.04_all/var/lib/gdm3/.config/systemd/user/drumroll.service
[Unit]
Description=Drumroll
Requires=pulseaudio.socket
After=systemd-user-sessions.service

[Service]
Type=simple
Restart=no
ExecStart=/usr/bin/paplay -p /usr/share/sounds/Yaru/stereo/system-ready.oga

[Install]
WantedBy=default.target
